# scripts/com.neelkanth.autodeploy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.neelkanth.autodeploy</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/sharadchaurasia/Documents/sharad/hotel-neelkanth-tracker/scripts/auto-deploy-full.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/sharadchaurasia/Documents/sharad/hotel-neelkanth-tracker/logs/autodeploy.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/sharadchaurasia/Documents/sharad/hotel-neelkanth-tracker/logs/autodeploy-error.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/sharadchaurasia/Documents/sharad/hotel-neelkanth-tracker</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
